# deploy/systemd/aialra-cleanup.timer
[Unit]
Description=Run AIALRA Localizer cleanup checks daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
